# root/etc/cron.d/pulledpork
30 2 * * * root /usr/bin/pulledpork -c /etc/snort/pulledpork.conf -l >/dev/null
