# root/etc/cron.d/pulledpork
30 2 * * * root /usr/bin/pulledpork -c /etc/pulledpork/pulledpork.conf -k -g -l >/dev/null; /etc/e-smith/events/actions/nethserver-pulledpork-apply >/dev/null
